# infra/ansible/roles/hermes/templates/hermes-config-webhook.service.j2
[Unit]
Description=GitHub webhook receiver for live Hermes config repository
Wants=network-online.target
After=network-online.target hermes-config-sync.service

[Service]
Type=simple
EnvironmentFile={{ hermes_config_sync_env_path }}
ExecStart=/usr/bin/python3 {{ hermes_config_webhook_path }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
