# infra/ansible/roles/hermes/templates/hermes-config-watch.service.j2
[Unit]
Description=Watch live Hermes config repository for local auto-improvement changes
Wants=network-online.target
After=network-online.target hermes-config-sync.service

[Service]
Type=simple
EnvironmentFile={{ hermes_config_sync_env_path }}
ExecStart=/bin/bash -lc 'while true; do inotifywait -r -e close_write,create,delete,move --exclude "(^|/)\\.git(/|$)" "$HERMES_CONFIG_DIR"; {{ hermes_config_git_sync_path }} --source inotify; done'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
