# infra/ansible/roles/hermes/templates/hermes-gateway.service.j2
[Unit]
Description=Hermes Agent Discord gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ hermes_user }}
Group={{ hermes_group }}
EnvironmentFile={{ hermes_env_path }}
WorkingDirectory={{ hermes_workspace }}
ExecStart={{ hermes_venv_path }}/bin/hermes gateway
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
