# infra/ansible/roles/docker_compose_project/templates/adguard-acme.service.j2
[Unit]
Description=Renew AdGuard Home certificate for Docker Compose DNS
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile={{ adguard_acme_env_path }}
ExecStart=/usr/local/bin/renew-adguard-cert
